# setup/runit/log/run
#!/bin/sh

exec chpst -unobody svlogd -tt ./main
